# supervisord.conf
[program:ws-server]
command=/usr/bin/php /home/admin/web/tracker.mycodesamples.com/public_html/bin/server.php
autostart=true
autorestart=true